# backend/supervisord.conf
[supervisord]
nodaemon=true

[program:php-fpm]
command=php-fpm
autostart=true
autorestart=true
stderr_logfile=/var/log/php-fpm.err.log
stdout_logfile=/var/log/php-fpm.out.log

[program:messenger-consumer-domain-events]
command=php /var/www/html/bin/console messenger:consume domain_events --time-limit=3600
autostart=true
autorestart=true
stderr_logfile=/var/log/messenger.err.log
stdout_logfile=/var/log/messenger.out.log

[program:messenger-consumer-notification-events]
command=php /var/www/html/bin/console messenger:consume notification_events --time-limit=3600
autostart=true
autorestart=true
stderr_logfile=/var/log/messenger-notification.err.log
stdout_logfile=/var/log/messenger-notification.out.log
